# Projects/dmm-manga-affiliate/Workflows/launchd/com.ai-brain.comment-poster.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-brain.comment-poster</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/tagishitakuya/Desktop/ClaudeProjects/AI-Brain/Projects/dmm-manga-affiliate/Workflows/launchd/run-comment-poster.sh</string>
    </array>

    <!-- 5分おき（300秒）-->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/ai-brain-comment-poster.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ai-brain-comment-poster.err</string>
</dict>
</plist>
