# Projects/dmm-manga-affiliate/Workflows/launchd/com.ai-brain.discord-reporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-brain.discord-reporter</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>source ~/.zprofile 2>/dev/null; /usr/bin/python3 /Users/tagishitakuya/Desktop/ClaudeProjects/AI-Brain/Projects/dmm-manga-affiliate/Workflows/discord-reporter.py</string>
    </array>
    <!-- 毎日 21:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/ai-brain-discord-reporter.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ai-brain-discord-reporter.err</string>
</dict>
</plist>
